# usb-decrypt@.service
[Service]
Type=simple
RuntimeMaxSec=660
ExecStart=/bin/sh -c "exec /opt/usb-auto-decrypt/decrypt_usb %I"
