# usb-decrypt@.service
[Service]
Type=forking
GuessMainPID=false
ExecStart=/bin/sh -c "exec /data/data/nas_stuff/decrypt_usb %I"
