# usb-decrypt@.service
[Service]
Type=simple
RuntimeMaxSec=660
ExecStart=/bin/sh -c "exec /data/data/nas_stuff/decrypt_usb %I"
