# usb-decrypt@.service
[Service]
Type=exec
RuntimeMaxSec=600
ExecStart=/bin/sh -c "exec /data/data/nas_stuff/decrypt_usb %I"
